# systemd/neon-image-sync.service
[Unit]
Description=Sync neon images

[Service]
Type=oneshot
ExecStart=/home/neon-image-syncd/sync
User=neon-image-sync
Group=neon-image-sync
