# systemd/neon-image-sync.timer
[Unit]
Description=Timer to sync neon images

[Timer]
OnCalendar=hourly
Persistent=true
